# data/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/php-fpm.log
stderr_logfile=/var/www/html/storage/logs/php-fpm.ERR.log
startsecs=5

[program:queue]
command=php /var/www/html/artisan queue:work --tries=3
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/queue.log
stderr_logfile=/var/www/html/storage/logs/queue.ERR.log
startsecs=5


[program:crontab]
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/crontab.log
stderr_logfile=/var/www/html/storage/logs/crontab.ERR.log
startsecs=5

